# raw_wireguard/rootfs/etc/services.d/wireguard/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: WireGuard
# Runs WireGuard
# ==============================================================================
declare interface

bashio::log.info "Starting WireGuard..."

# Get interface
interface="wg0"
if bashio::config.has_value "interface_name"; then
    interface=$(bashio::config "interface_name")
fi

# Run the WireGuard
wg-quick up "${interface}"

while true
do
    bashio::log.info "Requesting current status from WireGuard..."
    wg show
    sleep 30
done